# pki-update-crl.timer
[Unit]
Description=Update CRL periodically
Wants=network-online.target

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
